# english-bot.service
[Unit]
Description=English Telegram Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/english-prompts-bot
ExecStart=/opt/english-prompts-bot/venv/bin/python bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
